# integrations/telegram/systemd/claw-telegram-bridge.service
[Unit]
Description=Authenticated Telegram project bridge for Claw Code
After=network-online.target gemma-api-tunnel.service
Wants=network-online.target

[Service]
Type=simple
User=clawrun
Group=clawrun
WorkingDirectory=/home/clawrun/analise_storage
EnvironmentFile=/etc/claw-telegram-bridge.env
# Optional operator-managed credentials. Every variable loaded here is passed
# unchanged to Claw, child Agents, shell commands, plugins, and MCP processes.
EnvironmentFile=-/etc/claw-agent-credentials.env
ExecStart=/usr/bin/python3 /home/clawrun/analise_storage/claw-code-parity/integrations/telegram/claw_project_bridge.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillMode=mixed
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
# The '-' prefix prevents a first-install namespace failure before the bridge
# has created these writable directories. Deployment should still pre-create
# and chown them as documented.
ReadWritePaths=-/home/clawrun/analise_storage/state -/home/clawrun/analise_storage/telegram-projects -/home/clawrun/.claw

[Install]
WantedBy=multi-user.target
